# docker/services/php/conf.d/worker.ini
[program:worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/<project>/artisan queue:work --sleep=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=dev
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/<project>/public/worker.log
stopwaitsecs=3600